# speedtest/speedtest-zabbix.cron
32 * * * * zabbix /etc/zabbix/scripts/speedtest-zabbix.sh >/dev/null
